# systemd/backend_smm_bot.service
[Service]
User=root
WorkingDirectory=/home/rynoksmm/API
ExecStart=/usr/bin/python3 Backend.py
Restart=always

[Install]
WantedBy=multi-user.target
