# systemd/smm_bot.service
[Service]
User=root
WorkingDirectory=/home/rynoksmm
ExecStart=/usr/bin/python3 main.py
Restart=always

[Install]
WantedBy=multi-user.target
